# roles/pango-backup/templates/pango-backup.timer
[Unit]
Description=Schedule for backing up vms

[Timer]
OnCalendar="{{ backup_schedule }}"

[Install]
WantedBy=timers.target
